# s6/etc/services.d/default/finish
#!/usr/bin/execlineb -S0
#
# Stop the container if the service fails
#
# https://github.com/just-containers/s6-overlay#writing-an-optional-finish-script
# https://github.com/just-containers/s6-overlay/issues/87#issuecomment-142456561

if -n { s6-test $# -ne 0 }
if -n { s6-test ${1} -eq 256 }

s6-svscanctl -t /var/run/s6/services
